# include/supervisor/bw-server.conf
[program:bw-server]
command = th server.lua
stderr_logfile = /var/log/supervisor/bw-stderr.log
stdout_logfile = /var/log/supervisor/bw-stdout.log
startretries = 100
startsecs = 10
